# .config/systemd/nd/battery-check.target
[Unit]
Description=Check Battery Target
StopWhenUnneeded=yes
